# supervisor.ini
[program:wsgi]
command=.../gunicorn project.wsgi -b 0.0.0.0:8000 -w 9
directory=/home/admin/develop/django-template
autorestart=true
redirect_stderr=true
stdout_logfile=/home/admin/develop/log/django-template-wsgi.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info
stopsignal=KILL
stopasgroup=true

[program:asgi]
command=.../daphne project.asgi:application -b 0.0.0.0 -p 8001
directory=/home/admin/develop/django-template
autorestart=true
redirect_stderr=true
stdout_logfile=/home/admin/develop/log/django-template-asgi.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=50
stdout_capture_maxbytes=1MB
stdout_events_enabled=false
loglevel=info
stopsignal=KILL
stopasgroup=true
